# data/opkg/control/prerm
#!/bin/sh

systemctl stop unipi-control.service
systemctl disable unipi-control.service
systemctl --system daemon-reload
