# src/unipi_control/installer/etc/systemd/system/unipi-control.service
[Unit]
Description=Unipi Control
After=multi-user.target
Requires=unipitcp.service

[Service]
Type=simple
ExecStart=/usr/bin/unipi-control
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
